# systemd/www.socket
[Unit]
Description=www gunicorn socket

[Socket]
ListenStream=/run/www/socket

[Install]
WantedBy=sockets.target
